# components/070-rootfs/src/etc/inittab
::sysinit:/etc/rc
::respawn:-/sbin/shell
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/shutdown
